# configs/systemd/panic.service
[Unit]
Description=ZERO-DAY OS Panic Kill Switch (manual trigger only)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/panic

# This service should NOT be started at boot. It is triggered on-demand only.
# Remove WantedBy so 'systemctl enable' does not create a boot-time symlink.